# docker/etc/crontabs/nginx
0 8 * * * cd /var/www/app && ./cli cronjob >/dev/null 2>&1
0 */2 * * * cd /var/www/app && php pages/tools/staticsync.php
